# apps/api/tmp
# syntax=docker/dockerfile:1

ARG NODE_VERSION=25.8.2
ARG PNPM_VERSION=10.6.2

FROM node:25.8.2-alpine

RUN apk update
RUN apk add openrc
RUN apk add prometheus-node-exporter-openrc
RUN apk add nodejs npm
RUN node -v
RUN npm -v

ENV NODE_ENV development

RUN --mount=type=cache,target=/root/.npm \
    npm install -g pnpm@${PNPM_VERSION}

WORKDIR /usr/src/app

RUN --mount=type=bind,source=package.json,target=package.json \
    --mount=type=bind,source=pnpm-lock.yaml,target=pnpm-lock.yaml \
    --mount=type=cache,target=/root/.local/share/pnpm/store \
    pnpm install --frozen-lockfile

COPY . .

RUN wget https://github.com/prometheus/node_exporter/releases/download/v1.0.1/node_exporter-1.0.1.linux-amd64.tar.gz

RUN addgroup node_exporterg
RUN adduser -G node_exporterg -D -H node_exporter
RUN mkdir /etc/node_exporter
RUN chown node_exporter:node_exporterg /etc/node_exporter

RUN tar -xvf node_exporter-1.0.1.linux-amd64.tar.gz
RUN mv node_exporter-1.0.1.linux-amd64 node_exporter-files

RUN cp node_exporter-files/node_exporter /usr/local/bin/
RUN chown node_exporter:node_exporterg /usr/local/bin/node_exporter
RUN cat <<EOF > /etc/init.d/node_exporter
#!/sbin/openrc-run
name="node_exporter"
command="/usr/local/bin/node_exporter" # Path to your executable script
pidfile="/run/${RC_SVCNAME}.pid" # OpenRC variable for the service name
command_background="yes" # Use this if your command runs in the background

# Optional: Define dependencies (e.g., after the network is up)
depend() {
    after network
}
EOF

RUN chmod +x /etc/init.d/node_exporter
RUN chmod +x /usr/local/bin/node_exporter

EXPOSE 8000 9100 80 

CMD ["pnpm", "dev"]
